# crontab
0 4 * * * backup_jobs.sh > /logs/jobs.log 2>&1
0 4 * * * backup_departments.sh > /logs/departments.log 2>&1
0 4 * * * backup_hired_employees.sh > /logs/hired_employees.log 2>&1